# ansible/roles/zot_metrics/templates/zot-metrics.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!-- {{ ansible_managed }} -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>mcquack.eblume.zot-metrics</string>
	<key>ProgramArguments</key>
	<array>
		<string>{{ zot_metrics_script }}</string>
	</array>
	<key>StartInterval</key>
	<integer>{{ zot_metrics_interval }}</integer>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>{{ zot_metrics_log_dir }}/mcquack.zot-metrics.err.log</string>
	<key>StandardOutPath</key>
	<string>{{ zot_metrics_log_dir }}/mcquack.zot-metrics.out.log</string>
</dict>
</plist>
